# contrib/ed2k-server.service
# Example systemd unit for ed2k-server.
# Copy to /etc/systemd/system/ed2k-server.service on the VPS, edit paths
# as needed, then: systemctl daemon-reload && systemctl restart ed2k-server
#
# Key settings vs the default ulimit:
#   LimitNOFILE        — raise the file-descriptor limit. The default 1024 is
#                        fatally low. Fd usage tracks CONCURRENT CLIENTS, not the
#                        file index (the index is in RAM, costs zero fds). Budget
#                        ~1 fd per client + ~20 fixed (listeners/admin/log) +
#                        a margin for transient HighID-probe / gossip sockets.
#                        Rule of thumb: ~2 x max_clients. See README.
#   Restart=on-failure — auto-restart if the process exits unexpectedly.
#   Environment=RUST_LOG — log verbosity. NOTE: stdout/stderr are sent to
#                        /dev/null below (logs are OFF by default — they grow
#                        very fast). See the README before enabling them.

[Unit]
Description=ed2k index server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/ed2k-server --config /etc/ed2k-server/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

# Logs are DISABLED by default: they grow very fast on a busy server. To enable
# logging, change these to "journal" (view with: journalctl -u ed2k-server -f).
StandardOutput=null
StandardError=null

# File-descriptor limit. Each TCP client costs ~1 fd; the server also keeps
# the TCP + 5 UDP listeners, the admin TCP listener, log handles, and transient
# outbound sockets (HighID probes, gossip). The file index is in RAM and costs
# no fds. 65536 suits the default max_clients=1000 with huge margin. For a large
# server raise it to ~2 x max_clients (e.g. 131072 for max_clients=50000).
LimitNOFILE=65536

# Optional log level. Default is info; uncomment for verbose tracing.
# (Has no visible effect while StandardOutput/Error are "null" above.)
# Environment=RUST_LOG=ed2k_server=debug
Environment=RUST_LOG=ed2k_server=info

[Install]
WantedBy=multi-user.target
